# example_terraform/modules/launch_template/standalone/app.service
[Unit]
Description=Systemd Test Daemon

[Service]
Type=simple
ExecStart=/app
Restart=on-failure
